# dotfiles/.config/systemd/user/fix-steam-shortcuts.service
[Unit]
Description=Fix StartupWMClass on Steam shortcut .desktop files

[Service]
Type=oneshot
ExecStart=%h/.local/bin/fix-steam-shortcuts
